# crontab
* * * * * /app/cron-refresh.sh >> /dev/stdout 2>&1
